# mixi/src/bin/tool_batch_computebyday.cron
#! /bin/sh
cd $1
$2 -f tool_batch_computebyday.php